# amber.service
# systemd unit for "survival across reboots" — brings the Compose stack up on
# boot and down on shutdown. Deploys (rebuilds) are `bin/deploy`'s job, not this
# unit's; this just runs `docker compose up -d` against whatever image exists.
#
# Install on the box (one-time, and after editing this file):
#   sudo cp /home/amber/Amber/amber.service /etc/systemd/system/amber.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now amber.service
#
# WorkingDirectory must match the actual checkout path exactly — a case typo
# here (…/amber vs …/Amber) is what left this unit dead with 200/CHDIR.
#
# `-f compose.prod.yaml` is explicit because the repo root also ships
# `compose.yaml` (the desktop-developer install, mounts `spaces/example`).
# Without the flag, Docker Compose's default lookup would prefer that file
# and silently swap production onto the dev install on the next reboot.

[Unit]
Description=Amber Docker Compose stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/amber/Amber
ExecStart=/usr/bin/docker compose -f compose.prod.yaml up -d
ExecStop=/usr/bin/docker compose -f compose.prod.yaml down
User=amber
Group=amber

[Install]
WantedBy=multi-user.target
